# water/files/etc/rc.local
/SDC/gpsd/bin/gpsd -nN /dev/ttyS1
/SDC/SdcWatchDog &
exit 0
